# overlay/etc/cron.d/owncloud
# run cron.php for owncloud
*/15 * * * * www-data 	/usr/bin/php -f /var/www/owncloud/occ system:cron /dev/null 2>&1
